# docker/app/etc/services.d/serf/run
#!/usr/bin/with-contenv bash
pid=$$

sigterm() {
  kill -INT $pid
}

trap sigterm SIGTERM

/usr/local/bin/serf agent --config-file /etc/serf.conf --join "$SERF_JOIN" &
pid=$!
wait
