# docker/job/services.d/resque/run
#!/usr/bin/with-contenv bash
exec rake resque:work
